# crontab
* * * * * root /usr/bin/archive-wildfly-logs.sh >> /var/log/cron.log 2>&1
